# deploy/backtest-api.container
[Unit]
Description=Trading Backtest API
After=network-online.target

[Container]
Image=registry.starbluesolutions.net/backtest-api:latest
Pull=newer
ContainerName=backtest-api
PublishPort=127.0.0.1:8070:8070

Volume=/opt/trading/data:/app/data:ro,Z
Volume=/opt/trading/strategies:/app/strategies:Z
Volume=/opt/trading/backtest_results:/app/backtest_results:Z

[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
